# deploy/kid-buddy.service
[Unit]
Description=Reachy Kid Buddy -- real-time Gemini Live voice playmate
After=network-online.target

[Service]
Type=simple
User=pollen
WorkingDirectory=/home/pollen/everything-agent
Environment=HOME=/home/pollen
EnvironmentFile=/home/pollen/everything-agent/.env
ExecStart=/venvs/apps_venv/bin/python -u scripts/kid_buddy.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
